# soundbites.conf
description "uWSGI server for Soundbites: Staging"
start on runlevel [2345]
stop on runlevel [!2345]
respawn

script
    /usr/local/bin/uwsgi --ini /home/jbowers/apps/soundbites/app.ini
end script